# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:xvfb]
command=/usr/bin/Xvfb :0 -screen 0 1280x720x24
autostart=true
autorestart=true
stdout_logfile=/var/log/xvfb.log
stderr_logfile=/var/log/xvfb.err

[program:css]
command=/bin/bash -c "sleep 3 && /game/scripts/start_css.sh"
autostart=true
startsecs=0
stdout_logfile=/var/log/css.log
stderr_logfile=/var/log/css.err
environment=DISPLAY=:0,XDG_RUNTIME_DIR=/tmp/runtime-dir